# src/systemd/homesensors.service
[Unit]
Description=uWSGI instance to serve homesensors
After=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/home-sensors/src
Environment="PATH=/home/pi/home-sensors/src/homesensorsenv/bin"
ExecStart=/home/pi/home-sensors/src/homesensorsenv/bin/uwsgi --ini homesensors.ini

[Install]
WantedBy=multi-user.targe
